# mihomo/com.mihomo.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mihomo.proxy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mihomo/bin/mihomo</string>
        <string>-d</string>
        <string>/Users/mihomo/.config/mihomo</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/mihomo</string>
    
    <key>StandardOutPath</key>
    <string>/Users/mihomo/logs/mihomo.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mihomo/logs/mihomo.log</string>
    
    <key>UserName</key>
    <string>mihomo</string>
    
    <key>GroupName</key>
    <string>staff</string>
    
    <!-- 限制资源使用 -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
        <key>NumberOfProcesses</key>
        <integer>50</integer>
    </dict>
</dict>
</plist>