# mihomo/com.mihomo.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mihomo.proxy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mihomo/start-mihomo.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/mihomo</string>
    
    <key>StandardOutPath</key>
    <string>/Users/mihomo/logs/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mihomo/logs/stderr.log</string>
    
    <key>UserName</key>
    <string>mihomo</string>
    
    <key>GroupName</key>
    <string>staff</string>
    
    <!-- 限制资源使用 -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
        <key>NumberOfProcesses</key>
        <integer>10</integer>
    </dict>
    
    <!-- 网络访问权限 -->
    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SockServiceName</key>
            <string>7890</string>
            <key>SockType</key>
            <string>stream</string>
        </dict>
    </dict>
</dict>
</plist>